# templates/tftp.service.j2
[Unit]
Description=Tftp Server
Requires=tftp.socket
Documentation=man:in.tftpd

[Service]
ExecStart=/usr/sbin/in.tftpd -s {{ tftpboot_path }} -u {{ tftp_user }}
StandardInput=socket

[Install]
Also=tftp.socket
